# cron_refresh/crontab
* * * * * root /refresh.py; /script.sh;
