# cron/crontab
* * * * * root /refresh.py; /refresh.sh;
0 * * * * root /clean.py; /clean.sh
