# cron/crontab
* * * * * root /refresh.sh;
0 * * * * root /clean.sh
