# cron/crontab
1-59 * * * * root /request.sh refresh;
0 * * * * root /request.sh clean
