# cron/crontab
* * * * * root /request.sh refresh;
0 * * * * root /request.sh clean
